# rc.d/mosdns
#!/bin/sh

# PROVIDE: mosdns
# REQUIRE: SERVERS
# KEYWORD: shutdown

. /etc/rc.subr

name="mosdns"
rcvar="mosdns_enable"  # 允许通过 rc.conf 控制启用

desc="mosdns service"
pidfile="/var/run/${name}.pid"
logfile="/var/log/${name}.log"
procname="/usr/local/bin/mosdns"
command="/usr/sbin/daemon"
required_files="/usr/local/bin/mosdns /usr/local/etc/mosdns/config.yaml"
start_precmd="mosdns_prestart"
start_postcmd="mosdns_poststart"
stop_cmd="mosdns_stop"
status_cmd="mosdns_status"
extra_commands="reload"
reload_cmd="mosdns_reload"
mosdns_env="IS_DAEMON=1"  # 设置环境变量

# daemon 参数
command_args="-P ${pidfile} -f -o ${logfile} -t ${name}:daemon env ${mosdns_env} /usr/local/bin/mosdns start --as-service -d /usr/local/etc/mosdns/"

mosdns_prestart()
{
    touch "${logfile}"

    if [ -f "${pidfile}" ]; then
        if kill -0 "$(cat "${pidfile}")" 2>/dev/null; then
            echo "${name} is already running as pid $(cat "${pidfile}")."
            return 1
        fi
        rm -f "${pidfile}"
    fi

    if pgrep -f "/usr/local/bin/mosdns start --as-service -d /usr/local/etc/mosdns/" >/dev/null 2>&1; then
        echo "A mosdns process is already running; refusing to start another instance."
        return 1
    fi
}

mosdns_poststart()
{
    sleep 1

    if [ -f "${pidfile}" ]; then
        daemon_pid="$(cat "${pidfile}" 2>/dev/null)"
        if [ -n "${daemon_pid}" ] && kill -0 "${daemon_pid}" 2>/dev/null; then
            mosdns_pid="$(pgrep -P "${daemon_pid}" mosdns 2>/dev/null | head -n 1)"
            if [ -n "${mosdns_pid}" ] && kill -0 "${mosdns_pid}" 2>/dev/null; then
                echo "${name} started: daemon pid ${daemon_pid}, mosdns pid ${mosdns_pid}."
                return 0
            fi
        fi
    fi

    echo "${name} failed to start."
    return 1
}

mosdns_stop()
{
    stopped=0

    if [ -f "${pidfile}" ]; then
        daemon_pid="$(cat "${pidfile}" 2>/dev/null)"
        if [ -n "${daemon_pid}" ] && kill -0 "${daemon_pid}" 2>/dev/null; then
            mosdns_pids="$(pgrep -P "${daemon_pid}" mosdns 2>/dev/null)"

            if [ -n "${mosdns_pids}" ]; then
                echo "Stopping ${name} child process(es): ${mosdns_pids}"
                kill ${mosdns_pids} 2>/dev/null || true
            fi

            echo "Stopping ${name} wrapper process: ${daemon_pid}"
            kill "${daemon_pid}" 2>/dev/null || true
            stopped=1

            for _ in 1 2 3 4 5 6 7 8 9 10; do
                still_running=0

                if kill -0 "${daemon_pid}" 2>/dev/null; then
                    still_running=1
                fi

                for pid in ${mosdns_pids}; do
                    if kill -0 "${pid}" 2>/dev/null; then
                        still_running=1
                        break
                    fi
                done

                [ "${still_running}" -eq 0 ] && break
                sleep 1
            done

            if kill -0 "${daemon_pid}" 2>/dev/null; then
                kill -KILL "${daemon_pid}" 2>/dev/null || true
            fi

            for pid in ${mosdns_pids}; do
                if kill -0 "${pid}" 2>/dev/null; then
                    kill -KILL "${pid}" 2>/dev/null || true
                fi
            done
        fi
    fi

    extra_pids="$(pgrep -f "/usr/local/bin/mosdns start --as-service -d /usr/local/etc/mosdns/" 2>/dev/null)"
    if [ -n "${extra_pids}" ]; then
        echo "Stopping stray ${name} process(es): ${extra_pids}"
        kill ${extra_pids} 2>/dev/null || true
        stopped=1
        sleep 1
        for pid in ${extra_pids}; do
            if kill -0 "${pid}" 2>/dev/null; then
                kill -KILL "${pid}" 2>/dev/null || true
            fi
        done
    fi

    rm -f "${pidfile}"

    if pgrep -f "/usr/local/bin/mosdns start --as-service -d /usr/local/etc/mosdns/" >/dev/null 2>&1; then
        echo "${name} could not be stopped completely."
        return 1
    fi

    if [ "${stopped}" -eq 1 ]; then
        echo "${name} stopped."
    else
        echo "${name} is not running."
    fi

    return 0
}

mosdns_status()
{
    if [ -f "${pidfile}" ]; then
        daemon_pid="$(cat "${pidfile}" 2>/dev/null)"
        if [ -n "${daemon_pid}" ] && kill -0 "${daemon_pid}" 2>/dev/null; then
            mosdns_pid="$(pgrep -P "${daemon_pid}" mosdns 2>/dev/null | head -n 1)"
            if [ -n "${mosdns_pid}" ] && kill -0 "${mosdns_pid}" 2>/dev/null; then
                echo "${name} is running: daemon pid ${daemon_pid}, mosdns pid ${mosdns_pid}."
                return 0
            fi
            echo "${name} wrapper is running as pid ${daemon_pid}, but mosdns child process was not found."
            return 1
        fi
        echo "${name} pid file exists but process is not running."
        return 1
    fi

    mosdns_pid="$(pgrep -f "/usr/local/bin/mosdns start --as-service -d /usr/local/etc/mosdns/" 2>/dev/null | head -n 1)"
    if [ -n "${mosdns_pid}" ] && kill -0 "${mosdns_pid}" 2>/dev/null; then
        echo "${name} process is running as pid ${mosdns_pid}, but pidfile is missing."
        return 1
    fi

    echo "${name} is not running."
    return 1
}

mosdns_reload()
{
    if [ -f "${pidfile}" ]; then
        daemon_pid="$(cat "${pidfile}" 2>/dev/null)"
        if [ -n "${daemon_pid}" ] && kill -0 "${daemon_pid}" 2>/dev/null; then
            mosdns_pid="$(pgrep -P "${daemon_pid}" mosdns 2>/dev/null | head -n 1)"
            if [ -n "${mosdns_pid}" ] && kill -0 "${mosdns_pid}" 2>/dev/null; then
                kill -HUP "${mosdns_pid}" 2>/dev/null || return 1
                echo "${name} config reloaded."
                return 0
            fi
        fi
    fi

    echo "${name} is not running."
    return 1
}

# 加载用户定义的配置
load_rc_config "${name}"
: ${mosdns_enable:="NO"}  # 默认不启用

run_rc_command "$1"
